# apis_services/bressrain/bressrain.timer
[Unit]
Description=Bressrain scheduler

[Timer]
OnBootSec=0min
OnUnitActiveSec=5min

[Install]
WantedBy=multi-user.target